# scripts/services/pi-media-worker.service
[Unit]
Description=Pi Media Server - Processing Worker
After=network.target redis.target
Wants=redis.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-media-server
ExecStart=/usr/bin/node src/worker.js
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production

# Resource limits for video processing
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=90%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-media-worker

# Nice level for background processing
Nice=10

[Install]
WantedBy=multi-user.target